# rc.local
# ========================================================
# Sunberry
# Disable HDMI port
sudo /opt/vc/bin/tvservice -o
# ========================================================
